# deploy/systemd/einsatzueberwachung-backup.timer
[Unit]
Description=Taegliches Backup fuer Einsatzueberwachung

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=einsatzueberwachung-backup.service

[Install]
WantedBy=timers.target
