# deploy/systemd/einsatzueberwachung-server.service
[Unit]
Description=Einsatzueberwachung Server
After=network.target

[Service]
WorkingDirectory=/opt/einsatzueberwachung/server
ExecStart=/usr/bin/dotnet Einsatzueberwachung.Server.dll
Restart=always
RestartSec=10
SyslogIdentifier=einsatzueberwachung-server
User=einsatz
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://127.0.0.1:5000
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false
Environment=EINSATZUEBERWACHUNG_DATA_DIR=/opt/einsatzueberwachung/data
Environment=EINSATZUEBERWACHUNG_REPORT_DIR=/opt/einsatzueberwachung/data/berichte

[Install]
WantedBy=multi-user.target
