# deploy/systemd/einsatzueberwachung-backup.service
[Unit]
Description=Einsatzueberwachung Backup ausfuehren
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/einsatzueberwachung/scripts/backup.sh
User=root
Group=root
Environment=RETENTION_DAYS=14
